# jenkins-master/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
directory=/
command=/usr/sbin/sshd -D
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s_error.log

[program:jenkins]
command = /usr/bin/java -server -Djava.net.preferIPv4Stack=true -Dhudson.slaves.ChannelPinger.pingInterval=61 -Dhudson.remoting.Launcher.pingIntervalSec=3601 -Dhudson.remoting.Launcher.pingTimeoutSec=1000  -jar /usr/share/jenkins/jenkins.war --webroot=/var/cache/jenkins/war --httpPort=8080 --ajp13Port=-1
stdout_logfile = syslog
redirect_stderr = true
user = jenkins
